# deployment/vm/eth-collector.service
[Unit]
Description=Ethereum Real-Time Block Collector (Batch Mode: 5min flush)
After=network.target

[Service]
Type=simple
User=terry_eonlabs_com
WorkingDirectory=/home/terry_eonlabs_com/eth-collector
Environment="BATCH_INTERVAL_SECONDS=300"
ExecStart=/usr/bin/python3 -u /home/terry_eonlabs_com/eth-collector/realtime_collector.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eth-collector

[Install]
WantedBy=multi-user.target
